# root/etc/services.d/rescuetime/run
#!/usr/bin/execlineb -P

# Start rescuetime as jackson
# gosu jackson rescuetime
