# root/etc/services.d/rescuetime/run
#!/bin/sh

# Start rescuetime as jackson
exec gosu jackson rescuetime
