# root/etc/services.d/rescuetime/run
#!/bin/sh

# Start rescuetime as jackson
exec s6-setuidgid jackson rescuetime
