# root/etc/services.d/tmux/run
#!/usr/bin/execlineb -P

# Start tmux server as jackson
gosu jackson zsh -c 'tmux new-session -s desktop-environment zsh --login'
