# root/etc/services.d/tmux/run
#!/usr/bin/execlineb -P

# Start tmux server as jackson
exec gosu jackson \
  tmux new-session -s desktop-environment zsh --login
