# root/etc/services.d/tmux/run
#!/bin/sh

# Start tmux server as jackson
exec s6-setuidgid jackson \
  tmux new-session -s desktop-environment zsh --login
