# root/etc/services.d/webrelay/run
#!/usr/bin/execlineb -P

# Start webrelay tunnel as jackson
gosu jackson relay connect -s $HOSTNAME --crypto full http://localhost:8022
