# root/etc/services.d/webrelay/run
#!/usr/bin/execlineb -P

# Start webrelay tunnel as jackson
gosu jackson relay connect --name $HOSTNAME localhost:8022
