# docker/root/etc/services.d/transmission/run
#!/usr/bin/execlineb -P

# Start transmission daemon as jackson
gosu jackson transmission-daemon \
  --bind-address-ipv4 localhost \
  --config-dir $HOME/.config/transmission \
  --download-dir $HOME/torrents \
  --foreground \
  --no-auth \
  --rpc-bind-address localhost \
  --watch-dir $HOME/torrents/.watch
