# root/etc/services.d/keychain/run
#!/bin/sh

# Start ssh keychain as jackson
# exec gosu jackson keychain
