# root/etc/services.d/keychain/run
#!/usr/bin/execlineb -P

# Start ssh keychain as jackson
exec gosu jackson keychain
