# root/etc/services.d/keychain/run
#!/bin/sh

# Start ssh keychain as jackson
exec s6-setuidgid jackson keychain
