# root/etc/services.d/tmux-gotty/run
#!/usr/bin/execlineb -P

# Start gotty tmux session as jackson
gosu jackson tmux new-session -s gotty
